# supervisor.ini
[program:tbbd]
command=/app/tbbd/bin/gunicorn -b 0.0.0.0:8001 -w 4 webhook:app
directory=/app/tbbd/app
user=tbbd

[program:tbbd-image]
command=/app/tbbd/bin/python image.py
directory=/app/tbbd/app
user=tbbd
